# __config/crontab
* * * * * /opt/wexcommerce/__scripts/noip.sh >> /var/log/noip.log 2>&1
@reboot /opt/wexcommerce/__scripts/noip.sh >> /var/log/noip.log 2>&1
